# supervisor/supervisor.d/nginx.ini
[program:nginx]
command = nginx -g "daemon off;"
autostart = true
autorestart = true
stdout_logfile = /var/www/supervisor/log/nginx.log
stderr_logfile = /var/www/supervisor/log/nginx_err.log
