# supervisor/supervisor.d/node.ini
[program:node]
command = npm run prod
autostart = true
autorestart = true
stdout_logfile = /var/www/supervisor/log/node.log
stderr_logfile = /var/www/supervisor/log/node_err.log