# docker/supervisor/supervisor.d/node-dev.ini
[program:node-dev]
command = npm run dev
autostart = true
autorestart = true
stdout_logfile = /var/www/log/node.log
stderr_logfile = /var/www/log/node_err.log