# docker/supervisor/supervisor.d/node.ini
[program:node]
command = npm run build && npm run prod
autostart = true
autorestart = true
stdout_logfile = /var/www/log/node.log
stderr_logfile = /var/www/log/node_err.log